# systemd/foosnapper.service
[Unit]
Description=Automatic filesystem snapshooter
After=stratisd.service

[Service]
Type=oneshot
ExecStart=foosnapper
IOSchedulingClass=idle
