# systemd/foosnapper.service
[Unit]
Description=Automatic filesystem snapshooter
After=stratisd.service

[Service]
IOSchedulingClass=idle
ExecStart=/usr/bin/foosnapper
